# cron-job/crontab
* * * * * /bin/sh /home/rafaelfs/Área\ de\ Trabalho/monitor.sh >> /home/rafaelfs/Área\ de\ Trabalho/monitor.log 2>&1
# Single cron example